# xvfd.conf
description "Service for xvfd server"
author      "Dylan Grafmyre"

start on runlevel [2345]
stop on shutdown

exec /etc/init.d/xvfd
